# systemd/health-pipeline.service
[Unit]
Description=Health Pipeline Pull Analyze Notify
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/applewatch_openclaw_pipeline/openclaw_agent/run_health_pipeline_once.sh
WorkingDirectory=/root/applewatch_openclaw_pipeline
User=root
Group=root

[Install]
WantedBy=multi-user.target
